# ultimo-icondo_1.2.1-0_amd64/DEBIAN/postinst
#!/bin/bash
sudo rm -rf ~/.cache/mozilla/firefox/*
chmod +x ./home/administrator/ultimo/app/Ultimo.ICondo.WebApi
sudo systemctl enable ultimo.service
sudo systemctl start ultimo.service
