# ultimo-icondo_1.2.4-2_amd64/DEBIAN/postinst
#!/bin/bash
chmod +x ./home/administrator/ultimo/app/Ultimo.ICondo.WebApi
sudo rm -rf ~/.cache/mozilla/firefox/*
sudo systemctl enable ultimo.service
sudo systemctl start ultimo.service

